# infra/lightsail/clawos-telegram.service
[Unit]
Description=ClawOS Telegram Bot Adapter
Documentation=https://github.com/orestes-garcia-martinez/clawos
After=network.target
# Start after the skill worker if both are on the same instance.
Wants=clawos-worker.service

[Service]
Type=simple

# ── Security: non-root execution ──────────────────────────────────────────────
# Telegram adapter must never run as root.
User=clawos-admin
Group=clawos-admin

# ── Working directory and binary ──────────────────────────────────────────────
WorkingDirectory=/home/clawos-admin/clawos/apps/telegram
ExecStart=/usr/bin/node dist/index.js

# ── Environment ───────────────────────────────────────────────────────────────
# Secrets loaded from env file, never hardcoded here.
EnvironmentFile=/home/clawos-admin/clawos/apps/telegram/.env

# ── Restart policy ────────────────────────────────────────────────────────────
Restart=on-failure
RestartSec=5s
# Prevent restart loops from hammering the system.
StartLimitIntervalSec=60s
StartLimitBurst=3

# ── Hardening ─────────────────────────────────────────────────────────────────
NoNewPrivileges=true
ProtectSystem=strict
# /home/clawos-admin is readable for the deploy dir.
ReadWritePaths=/home/clawos-admin/clawos/apps/telegram
PrivateTmp=true

# ── Logging ───────────────────────────────────────────────────────────────────
# JSON audit entries go to stdout -> journald.
# View with: journalctl -u clawos-telegram -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clawos-telegram

[Install]
WantedBy=multi-user.target
